# Sample Data/MITRE ATT&CK/Layers/Linux.json
{"description": "Linux ATT&CK Matrix Coverage for Azure Sentinel", "name": "Linux Layer Json File for Azure Sentinel", "domain": "mitre-enterprise", "version": "3.0", "filters": {"stages": ["act"], "platforms": ["Linux"]}, "techniques": [{"score": 1, "techniqueID": "T1098", "metadata": [{"name": "GitHub User Grants Access and Other User Grants Access", "value": "GitHub_CL"}, {"name": "GitHub Update Permissions", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1078", "metadata": [{"name": "New internet-exposed SSH endpoints", "value": "Syslog"}, {"name": "IP with multiple failed Azure AD logins successfully logs in to Palo Alto VPN", "value": "CommonSecurityLog"}, {"name": "Cisco - firewall block but success logon to Azure AD", "value": "CommonSecurityLog"}, {"name": "Failed AzureAD logons but success logon to host", "value": "Syslog"}, {"name": "Multiple Password Reset by user", "value": "Syslog"}, {"name": "Failed host logons but success logon to AzureAD", "value": "Syslog"}, {"name": "Tracking Password Changes", "value": "Syslog"}, {"name": "GitHub Activites from a New Country", "value": "GitHub_CL"}, {"name": "User made Owner of multiple teams", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "T1136", "metadata": [{"name": "GitHub Inactive or New Account Access or Usage", "value": "GitHub_CL"}, {"name": "GitHub First Time Invite Member and Add Member to Repo", "value": "GitHub_CL"}, {"name": "External user added and removed in short timeframe", "value": "TeamsLogs_CL"}, {"name": "External user from a new organisation added", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "", "metadata": [{"name": "(Preview) TI map IP entity to SigninLogs", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Domain entity to PaloAlto", "value": "CommonSecurityLog"}, {"name": "(Preview) TI map Domain entity to PaloAlto", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Domain entity to SecurityAlert", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to AzureNetworkAnalytics_CL (NSG Flow Logs)", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to AWSCloudTrail", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map File Hash to Security Event", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Email entity to SecurityAlert", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Email entity to CommonSecurityLog", "value": "CommonSecurityLog"}, {"name": "(Preview) TI map Email entity to CommonSecurityLog", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Domain entity to CommonSecurityLog", "value": "ThreatIntelligenceIndicator"}, {"name": "Microsoft COVID-19 file hash indicator matches", "value": "CommonSecurityLog"}, {"name": "(Preview) TI map IP entity to W3CIISLog", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map URL entity to Syslog data", "value": "Syslog"}, {"name": "(Preview) TI map URL entity to Syslog data", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to VMConnection", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to VMConnection", "value": "VMConnection"}, {"name": "(Preview) TI map Email entity to OfficeActivity", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map URL entity to SecurityAlert data", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Email entity to SecurityEvent", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to DnsEvents", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to DnsEvents", "value": "DnsEvents"}, {"name": "(Preview) TI map File Hash to CommonSecurityLog Event", "value": "CommonSecurityLog"}, {"name": "(Preview) TI map File Hash to CommonSecurityLog Event", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Domain entity to Syslog", "value": "Syslog"}, {"name": "(Preview) TI map Domain entity to Syslog", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map URL entity to AuditLogs", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to OfficeActivity", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to WireData", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map IP entity to WireData", "value": "WireData"}, {"name": "(Preview) TI map IP entity to AzureActivity", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map URL entity to OfficeActivity data", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Email entity to AzureActivity", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Domain entity to DnsEvent", "value": "DnsEvents"}, {"name": "(Preview) TI map Domain entity to DnsEvent", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map URL entity to PaloAlto data", "value": "CommonSecurityLog"}, {"name": "(Preview) TI map URL entity to PaloAlto data", "value": "ThreatIntelligenceIndicator"}, {"name": "(Preview) TI map Email entity to SigninLogs", "value": "ThreatIntelligenceIndicator"}, {"name": "Known IRIDIUM IP", "value": "DnsEvents"}, {"name": "Known IRIDIUM IP", "value": "VMConnection"}, {"name": "Known IRIDIUM IP", "value": "CommonSecurityLog"}, {"name": "Known IRIDIUM IP", "value": "WireData"}, {"name": "THALLIUM domains included in DCU takedown", "value": "DnsEvents"}, {"name": "THALLIUM domains included in DCU takedown", "value": "VMConnection"}, {"name": "THALLIUM domains included in DCU takedown", "value": "CommonSecurityLog"}, {"name": "Known GALLIUM domains and hashes", "value": "DnsEvents"}, {"name": "Known GALLIUM domains and hashes", "value": "VMConnection"}, {"name": "Known GALLIUM domains and hashes", "value": "CommonSecurityLog"}, {"name": "TI map IP entity to GitHub_CL", "value": "ThreatIntelligenceIndicator"}, {"name": "Threats detected by Eset", "value": "eset_CL"}, {"name": "Web sites blocked by Eset", "value": "eset_CL"}, {"name": "Preview - TI map File entity to Syslog Event", "value": "Syslog"}, {"name": "Preview - TI map File entity to Syslog Event", "value": "ThreatIntelligenceIndicator"}, {"name": "Preview - TI map File entity to Security Event", "value": "ThreatIntelligenceIndicator"}, {"name": "Preview - TI map File entity to VMConnection Event", "value": "VMConnection"}, {"name": "Preview - TI map File entity to VMConnection Event", "value": "ThreatIntelligenceIndicator"}, {"name": "Preview - DNS Events that match threat intelligence", "value": "DnsEvents"}, {"name": "Preview - DNS Events that match threat intelligence", "value": "ThreatIntelligenceIndicator"}, {"name": "Preview - TI map File entity to OfficeActivity Event", "value": "ThreatIntelligenceIndicator"}, {"name": "Preview - TI map File entity to WireData Event", "value": "WireData"}, {"name": "Preview - TI map File entity to WireData Event", "value": "ThreatIntelligenceIndicator"}, {"name": "Web shell command alert enrichment", "value": "SecurityAlert"}, {"name": "Web shell file alert enrichment", "value": "SecurityAlert"}, {"name": "GitHub Security Vulnerability in Repository", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1100", "metadata": [{"name": "Malicious web application requests linked with MDATP alerts", "value": "SecurityAlert"}, {"name": "GitHub OAuth App Restrictions Disabled", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1190", "metadata": [{"name": "New High Severity Vulnerability Detected Across Multiple Hosts", "value": "QualysHostDetection_CL"}, {"name": "High Number of Urgent Vulnerabilities Detected", "value": "QualysHostDetection_CL"}, {"name": "Suspected use of attack framework followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1110", "metadata": [{"name": "New internet-exposed SSH endpoints", "value": "Syslog"}, {"name": "SSH - Potential Brute Force", "value": "Syslog"}, {"name": "Failed logon attempts in authpriv", "value": "Syslog"}, {"name": "IP with multiple failed Azure AD logins successfully logs in to Palo Alto VPN", "value": "CommonSecurityLog"}, {"name": "Failed AzureAD logons but success logon to host", "value": "Syslog"}, {"name": "Multiple Password Reset by user", "value": "Syslog"}, {"name": "Failed host logons but success logon to AzureAD", "value": "Syslog"}, {"name": "PulseConnectSecure - Potential Brute Force Attempts", "value": "Syslog"}, {"name": "PulseConnectSecure - Large Number of Distinct Failed User Logins", "value": "Syslog"}, {"name": "ClientDeniedAccess", "value": "Syslog"}, {"name": "Excessive Failed Authentication from Invalid Inputs", "value": "Syslog"}, {"name": "Disabled accounts using Squid proxy", "value": "Syslog"}, {"name": "Tracking Password Changes", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1059", "metadata": [{"name": "Rare process running on a Linux host", "value": "Syslog"}, {"name": "Crypto currency miners EXECVE", "value": "Syslog"}, {"name": "Editing Linux scheduled tasks through Crontab", "value": "Syslog"}, {"name": "Linux scheduled task Aggregation", "value": "Syslog"}, {"name": "PowerShell made a suspicious network connection, followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1485", "metadata": [{"name": "GitHub First Time Repo Delete", "value": "GitHub_CL"}, {"name": "GitHub Mass Deletion of repos or projects", "value": "GitHub_CL"}, {"name": "Multiple Teams deleted by a single user", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "T1090", "metadata": [{"name": "Squid proxy events for ToR proxies", "value": "Syslog"}, {"name": "User Accessed Suspicious URL Categories", "value": "Syslog"}, {"name": "Excessive Denied Proxy Traffic", "value": "Syslog"}, {"name": "Network request to TOR anonymization service followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1008", "metadata": [{"name": "Squid proxy events for ToR proxies", "value": "Syslog"}, {"name": "Excessive NXDOMAIN DNS Queries", "value": "Syslog"}, {"name": "Potential DGA detected", "value": "DnsEvents"}, {"name": "Squid commonly abused TLDs", "value": "Syslog"}, {"name": "Cobalt Strike DNS Beaconing", "value": "DnsEvents"}, {"name": "Cobalt Strike DNS Beaconing", "value": "VMConnection"}, {"name": "DNS - domain anomalous lookup increase", "value": "DnsEvents"}, {"name": "Abnormally long DNS URI queries", "value": "DnsEvents"}, {"name": "DNS Full Name anomalous lookup increase", "value": "DnsEvents"}, {"name": "DNS lookups for commonly abused TLDs", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1102", "metadata": [{"name": "Squid proxy events related to mining pools", "value": "Syslog"}, {"name": "Request for single resource on domain", "value": "CommonSecurityLog"}]}, {"score": 1, "techniqueID": "T1483", "metadata": [{"name": "Excessive NXDOMAIN DNS Queries", "value": "Syslog"}, {"name": "Potential DGA detected", "value": "DnsEvents"}, {"name": "Possible contact with a domain generated by a DGA", "value": "CommonSecurityLog"}, {"name": "Squid commonly abused TLDs", "value": "Syslog"}, {"name": "Cobalt Strike DNS Beaconing", "value": "DnsEvents"}, {"name": "Cobalt Strike DNS Beaconing", "value": "VMConnection"}, {"name": "DNS - domain anomalous lookup increase", "value": "DnsEvents"}, {"name": "Abnormally long DNS URI queries", "value": "DnsEvents"}, {"name": "DNS Full Name anomalous lookup increase", "value": "DnsEvents"}, {"name": "DNS lookups for commonly abused TLDs", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1200", "metadata": [{"name": "Potential DHCP Starvation Attack", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1043", "metadata": [{"name": "Known Phosphorus group domains/IP", "value": "DnsEvents"}, {"name": "Known Phosphorus group domains/IP", "value": "VMConnection"}, {"name": "Known Phosphorus group domains/IP", "value": "CommonSecurityLog"}, {"name": "Known Strontium group domains", "value": "DnsEvents"}, {"name": "Known Strontium group domains", "value": "VMConnection"}, {"name": "Known Strontium group domains", "value": "CommonSecurityLog"}, {"name": "Palo Alto - potential beaconing detected", "value": "CommonSecurityLog"}, {"name": "Fortinet - Beacon pattern detected", "value": "CommonSecurityLog"}, {"name": "Squid data volume timeseries anomalies", "value": "Syslog"}, {"name": "Retrospective hunt for STRONTIUM IP IOCs", "value": "CommonSecurityLog"}, {"name": "RareDNSLookupWithDataTransfer", "value": "DnsEvents"}, {"name": "RareDNSLookupWithDataTransfer", "value": "CommonSecurityLog"}, {"name": "RareDNSLookupWithDataTransfer", "value": "WireData"}, {"name": "RareDNSLookupWithDataTransfer", "value": "VMConnection"}, {"name": "Detect beacon like pattern based on repetitive time intervals in Wire Data Traffic", "value": "WireData"}]}, {"score": 1, "techniqueID": "T1071", "metadata": [{"name": "Request for single resource on domain", "value": "CommonSecurityLog"}, {"name": "Outbound connection to IP with a history of unauthorized access attempts followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1203", "metadata": [{"name": "Suspected use of attack framework followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1204", "metadata": [{"name": "Network endpoint to host executable correlation", "value": "CommonSecurityLog"}, {"name": "Known Malware Detected", "value": "CarbonBlackEvents_CL"}]}, {"score": 1, "techniqueID": "T1030", "metadata": [{"name": "Time series anomaly for data size transferred to public internet", "value": "CommonSecurityLog"}, {"name": "Time series anomaly for data size transferred to public internet", "value": "VMConnection"}, {"name": "Time series anomaly detection for total volume of traffic", "value": "CommonSecurityLog"}, {"name": "Squid data volume timeseries anomalies", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1089", "metadata": [{"name": "GitHub Two Factor Auth Disable", "value": "GitHub_CL"}, {"name": "GitHub Update Permissions", "value": "GitHub_CL"}, {"name": "GitHub OAuth App Restrictions Disabled", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1489", "metadata": [{"name": "Multiple Teams deleted by a single user", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "T1213", "metadata": [{"name": "GitHub Repo Clone - Time Series Anomly", "value": "GitHub_CL"}, {"name": "GitHub Repo switched from private to public", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1046", "metadata": [{"name": "Rare client observed with high reverse DNS lookup count", "value": "DnsEvents"}, {"name": "Palo Alto - possible internal to external port scanning", "value": "CommonSecurityLog"}, {"name": "Cisco ASA - average attack detection rate increase", "value": "CommonSecurityLog"}, {"name": "Cisco ASA - threat detection message fired", "value": "CommonSecurityLog"}, {"name": "Port Scan Detected", "value": "Syslog"}, {"name": "Squid malformed requests", "value": "Syslog"}, {"name": "High reverse DNS count by host", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1035", "metadata": [{"name": "DNS events related to mining pools", "value": "DnsEvents"}, {"name": "DNS Domains linked to WannaCry ransomware campaign", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1496", "metadata": [{"name": "DNS events related to mining pools", "value": "DnsEvents"}, {"name": "DNS Domains linked to WannaCry ransomware campaign", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1048", "metadata": [{"name": "DNS events related to ToR proxies", "value": "DnsEvents"}, {"name": "RareDNSLookupWithDataTransfer", "value": "DnsEvents"}, {"name": "RareDNSLookupWithDataTransfer", "value": "CommonSecurityLog"}, {"name": "RareDNSLookupWithDataTransfer", "value": "WireData"}, {"name": "RareDNSLookupWithDataTransfer", "value": "VMConnection"}, {"name": "DNS - domain anomalous lookup increase", "value": "DnsEvents"}, {"name": "Abnormally long DNS URI queries", "value": "DnsEvents"}, {"name": "DNS Full Name anomalous lookup increase", "value": "DnsEvents"}, {"name": "DNS lookups for commonly abused TLDs", "value": "DnsEvents"}]}, {"score": 1, "techniqueID": "T1065", "metadata": [{"name": "Palo Alto - potential beaconing detected", "value": "CommonSecurityLog"}, {"name": "Fortinet - Beacon pattern detected", "value": "CommonSecurityLog"}, {"name": "Detect beacon like pattern based on repetitive time intervals in Wire Data Traffic", "value": "WireData"}]}, {"score": 1, "techniqueID": "T1498", "metadata": [{"name": "Cisco ASA - average attack detection rate increase", "value": "CommonSecurityLog"}, {"name": "Cisco ASA - threat detection message fired", "value": "CommonSecurityLog"}]}, {"score": 1, "techniqueID": "T1210", "metadata": [{"name": "Critical Threat Detected", "value": "CarbonBlackNotifications_CL"}, {"name": "Suspected use of attack framework followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1499", "metadata": [{"name": "Excessive Amount of Denied Connections from a Single Source", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1053", "metadata": [{"name": "Rare process running on a Linux host", "value": "Syslog"}, {"name": "Crypto currency miners EXECVE", "value": "Syslog"}, {"name": "Editing Linux scheduled tasks through Crontab", "value": "Syslog"}, {"name": "Linux scheduled task Aggregation", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1037", "metadata": [{"name": "Rare process running on a Linux host", "value": "Syslog"}, {"name": "Editing Linux scheduled tasks through Crontab", "value": "Syslog"}, {"name": "Linux scheduled task Aggregation", "value": "Syslog"}]}, {"score": 1, "techniqueID": "T1176", "metadata": [{"name": "Previously unseen bot or applicaiton added to Teams", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "T1119", "metadata": [{"name": "Previously unseen bot or applicaiton added to Teams", "value": "TeamsLogs_CL"}]}, {"score": 1, "techniqueID": "T10178", "metadata": [{"name": "GitHub User Grants Access and Other User Grants Access", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1537", "metadata": [{"name": "GitHub OAuth App Restrictions Disabled", "value": "GitHub_CL"}]}, {"score": 1, "techniqueID": "T1047", "metadata": [{"name": "Suspicious remote WMI execution followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1573", "metadata": [{"name": "Network request to TOR anonymization service followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}, {"score": 1, "techniqueID": "T1068", "metadata": [{"name": "Suspected use of attack framework followed by anomalous traffic flagged by Palo Alto Networks firewall", "value": "SecurityAlert"}]}], "gradient": {"colors": ["#ffffff", "##8ec843"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "Techniques researched", "color": "#66fff3"}]}